# iBack/.config/supervisor/supervisor.conf
[program:uwsgi]
command = /root/.pyenv/versions/app_env/bin/uwsgi --ini /srv/app/.config/uwsgi/deploy.ini

[program:nginx]
command = nginx